# .local/bin/mail.cron
#!/bin/bash

PID=$(pgrep offlineimap)
LOCK=/var/lock/mail

# offlineimap freezes up a lot.  kill if previous job is still running
[[ -n "$PID" ]] && kill $PID

# sub-shell below is locked to prevent concurrency
(
	# bail early if lock is held
	flock -n 9 || exit 1

	/usr/bin/offlineimap -o -a baylibre -u basic | /usr/bin/logger
) 9>$LOCK

exit 0

# run this script periodically with the following cron job
# */15 * * * *  ~/.local/bin/mail.cron
