# .local/bin/mail.cron
#!/usr/local/bin/fish
#!/bin/bash

#PID=$(pgrep -f offlineimap)

# offlineimap freezes up a lot. kill if previous job is still running
#[[ -n "$PID" ]] && kill -9 $PID

/usr/local/bin/offlineimap -o -a baylibre -u basic | /usr/bin/logger

#exit 0

# run this script periodically with the following cron job
# */15 * * * *  ~/.local/bin/mail.cron
